# systemd/aienrich.service
[Unit]
Description=AI Enrich — enrich article URLs from Google Sheet to Notion
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=fastlypr
WorkingDirectory=/home/fastlypr/aienrich
# Common locations for python3 and the codex CLI. Adjust if your codex
# binary lives elsewhere (run `which codex` to confirm).
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/home/fastlypr/.local/bin:/home/fastlypr/.npm-global/bin:/home/fastlypr/.nvm/versions/node/current/bin"
ExecStart=/usr/bin/python3 -u enrich_urls_codex.py
# Output goes to journald — view with `journalctl -u aienrich -f`
StandardOutput=journal
StandardError=journal
